# admin/misc/lsws.rc.gentoo
#!/sbin/runscript

depend()
{
	need net
	use mysql dns netmount postgres
	after sshd
}

start()
{
	ebegin "Starting LiteSpeed Web Server"
	"/usr/local/lsws/bin/lswsctrl" start 2>/dev/null 1>&2
	eend $?
}

stop()
{
	ebegin "Stopping LiteSped Web Server"
	"/usr/local/lsws/bin/lswsctrl" stop 2>/dev/null 1>&2
	eend $?
}

restart()
{
	ebegin "Restarting LiteSpeed Web Server"
	"/usr/local/lsws/bin/lswsctrl" restart 2>/dev/null 1>&2
	eend $?
}
